# iac/ansible/roles/roles/samba4/templates/samba.service.j2
[Unit]
Description= Samba 4 AD Daemon
After=syslog.target
After=network.target

[Service]
Type=forking
PIDFile=/run/samba/samba.pid
ExecStart=/usr/sbin/samba

[Install]
WantedBy=multi-user.target